# configuration/etc/systemd/user/monitor-hotplug.service
[Unit]
Description=Handle monitor hotplug events

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-hotplug-bb
